# managed/cron/crontab
0 12 * * * DOTFILES_DIR/managed/cron/sync_homebrew_packages.sh  # MANAGED_BY_DOTFILES
0 13 * * * DOTFILES_DIR/managed/cron/sync_macos_applications.sh  # MANAGED_BY_DOTFILES
